# deploy/systemd/carelink-bridge.service
[Unit]
Description=carelink-bridge (Medtronic CareLink -> Nightscout)
Documentation=https://github.com/NovaLux12/carelink-bridge
After=network-online.target
Wants=network-online.target

# Don't auto-restart too aggressively — CareLink has rate limits and
# rapid restart loops can hammer their auth endpoint. 30s gives us
# breathing room while still recovering quickly from real failures.
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
# %h is the home directory of the user running this --user service, so the
# unit is portable across machines and accounts (install into
# ~/carelink-bridge as the deploy README describes) with no editing.
WorkingDirectory=%h/carelink-bridge
ExecStart=/usr/bin/node %h/carelink-bridge/dist/main.js
Restart=on-failure
RestartSec=30
TimeoutStopSec=15

# Env file holds the secrets (CARELINK_*, API_SECRET). The "-" prefix
# means missing-file-is-OK at startup, so the unit still loads even
# before the .env file is created post-deploy. The bridge itself
# will exit with a clear error if vars are missing at runtime.
EnvironmentFile=-%h/carelink-bridge/.env

# Explicit PATH — node lives at /usr/bin/node; node_modules .bin
# not on PATH for the service user.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:%h/.npm-global/bin
Environment=NODE_ENV=production

# Hardening (defence in depth, even for a personal service):
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
# Directory-level grant, not per-file: the bridge rewrites logindata.json
# on token refresh and unlinks it when the refresh token is dead, and
# unlinking a per-file ReadWritePaths bind mount fails with EBUSY (which
# client.ts swallows, so recovery would silently not happen). .env is
# read-only for the bridge and needs no write grant at all. See issue #16.
ReadWritePaths=%h/carelink-bridge

# Allow only outbound network (CareLink + Nightscout are HTTPS, so
# restricting to HTTPS ports keeps the surface tight). CareLink uses
# port 443. The bridge doesn't need inbound.
RestrictAddressFamilies=AF_INET AF_INET6
SystemCallArchitectures=native

StandardOutput=journal
StandardError=journal
SyslogIdentifier=carelink-bridge

[Install]
WantedBy=default.target
